# scripts/vibe-tanks-admin.service
[Unit]
Description=Vibe Tanks admin dashboard sidecar
After=network-online.target vibe-tanks.service
Wants=network-online.target
PartOf=vibe-tanks.service

[Service]
Type=simple
User=blin2h
Group=blin2h
WorkingDirectory=/home/blin2h/deploy/Vibe-Tanks/admin
ExecStart=/home/blin2h/.nvm/versions/node/v22.22.0/bin/node /home/blin2h/deploy/Vibe-Tanks/admin/node_modules/ts-node/dist/bin.js src/index.ts
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
Environment=GAME_INTERNAL_URL=http://127.0.0.1:3010

[Install]
WantedBy=multi-user.target
